# tests/containers/debian-snmp/supervisord.conf
[supervisord]
nodaemon=true
 
[program:sshd]
command=/usr/sbin/sshd -D

[program:snmpsim]
command=snmpsimd.py --daemonize --process-user=snmpsim --process-group=snmpsim --agent-udpv4-endpoint=0.0.0.0:161 --data-dir=/opt

[program:show]
command=/sbin/show init_context
